# files/tcp-dns-proxy
#!/sbin/runscript

APP_PATH="/home/ronald/program/Tcp-DNS-proxy/tcpdns.py"
PID_PATH="/var/run/tcpdns.pid"
start(){
    # Start the daemon 
    ebegin "Starting TCP DNS Proxy"
    	start-stop-daemon --background --name tcpDnsProxy --start --make-pidfile --pidfile $PID_PATH --exec python -- $APP_PATH
    eend $?
}
stop(){
    ebegin "Stopping TCP DNS Proxy"
    	start-stop-daemon --name tcpDnsProxy --stop python -- $APP_PATH
    eend $?
}
